# roles/cart/files/cart.service
[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis-dev.roboshop.blog
Environment=CATALOGUE_HOST=catalogue-dev.roboshop.blog
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target